# snapshots/20260317_003820/systemd-user/dgx-commands-api.service
[Unit]
Description=Commands API (port 8092) — Mission Control backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pmello/.openclaw/hub
ExecStart=/usr/bin/python3 commands-api.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
