# snapshots/20260301_154357/systemd-user/max-files-api.service
[Unit]
Description=MAX Files API (port 8093)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pmello/.openclaw/workspace-max/coding-tests/04-system-dashboard/files-api.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
