# snapshots/20260309_220402/systemd-user/max-commands-api.service
[Unit]
Description=MAX Commands API (port 8092)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pmello/.openclaw/workspace-max/coding-tests/04-system-dashboard/commands-api.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
