# snapshots/20260301_094356/systemd-user/jess-q35.service
[Unit]
Description=Jess Q35 llama.cpp server
After=network.target

[Service]
Type=simple
ExecStart=/home/pmello/llama.cpp-q35/build-cpu/bin/llama-server -m /home/pmello/models/qwen3.5-35b-a3b-q8/Qwen3.5-35B-A3B-Q8_0.gguf --host 127.0.0.1 --port 18080 --ctx-size 131072 --threads 16 --jinja --reasoning-format none --reasoning-budget 0
Restart=on-failure
RestartSec=3
WorkingDirectory=/home/pmello/.openclaw/workspace-q35

[Install]
WantedBy=default.target
