# snapshots/20260301_094356/systemd-user/max-web-gallery.service
[Unit]
Description=MAX Web Gallery (port 8090)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pmello/.openclaw/workspace-max/coding-tests
ExecStart=/usr/bin/python3 -m http.server 8090 --bind 0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
