# snapshots/20260302_034358/systemd-user/jess-keeper.service
[Unit]
Description=Jess keeper agent (memory scoring + RAG records)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/pmello/.openclaw/workspace-q35/scripts/jess_keeper.py
