# snapshots/20260307_012519/systemd-user/dgx-youtube-understand.service
[Unit]
Description=YouTube Understanding API
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pmello/.openclaw/workspace-main/tools/youtube-understand
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 /home/pmello/.openclaw/workspace-main/tools/youtube-understand/server.py
Restart=on-failure
RestartSec=5
Environment=HOME=/home/pmello
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target