# snapshots/20260307_210018/systemd-user/jess-keeper.timer
[Unit]
Description=Run Jess keeper every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
Persistent=true

[Install]
WantedBy=timers.target
